# .config/systemd/user/wireguard_watchdog.timer
[Unit]
Description=Periodically check if Wireguard is up

[Timer]
OnCalendar=*:*:0/30

[Install]
WantedBy=timers.target
